# maintenance/crontab
# Email Setting
MAILTO=<email@email.com>

# Jobs
0 */2 * * * /backup.sh 2>&1
